# slurm/Exposed_REST_API/scaleout/postfix.service
[Service]
ExecCondition=+bash -c 'exec test $(hostname -s) = "login"'
ExecStartPre=
