# slurm/Exposed_REST_API/scaleout/slurmrestd.service
[Service]
Environment="SLURM_JWT="
ExecCondition=bash -c 'test $(hostname -s) = "rest"'
ExecStartPre=+/usr/local/bin/slurmrestd.startup.sh
