# slurm/Exposed_REST_API/scaleout/slurmd.slice
[Unit]
Description=slurm slice
Before=slices.target
[Slice]
CPUAccounting=true
MemoryAccounting=true
Delegate=yes
